# copilot/plugins/create_ap/supervisor.conf
[program:create_ap]
command =   /home/www/copilot/copilot/plugins/create_ap/start
directory = /home/www/copilot/copilot/plugins/create_ap
user = root
autostart=true
autorestart=true
