# copilot/plugins/dns/supervisor.conf

[program:dnschef]
command =   /home/www/copilot/copilot/plugins/dns/start
directory = /home/www/copilot/copilot/plugins/dns
user = root
autostart=true
autorestart=true
