# docs/example_package/supervisor.conf
[program:example_package]
command =   /home/www/copilot/docs/example_package/start
user = root
autostart=true
autorestart=true